# trae_ai_base44/launchd/com.traeai.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.traeai.watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>scripts/watchdog.sh</string>
  </array>
  <key>WorkingDirectory</key><string>__REPO_ROOT__</string>
  <key>StandardOutPath</key><string>logs/watchdog.out.log</string>
  <key>StandardErrorPath</key><string>logs/watchdog.err.log</string>
  <key>StartInterval</key><integer>300</integer>
  <key>RunAtLoad</key><true/>
</dict>
</plist>